# systemd/system/snap-gnome\x2d3\x2d34\x2d1804-36.mount
[Unit]
Description=Mount unit for gnome-3-34-1804, revision 36
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/gnome-3-34-1804_36.snap
Where=/snap/gnome-3-34-1804/36
Type=squashfs
Options=nodev,ro,x-gdu.hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
